# templates/etc/init/vulcand.conf.j2
description "vulcand"

start on runlevel [2345]

exec {{ vulcand_install_path }}/vulcand{% for key, val in vulcand_args.iteritems() %} -{{ key }}={{ val }}{% endfor %}
